# dockerfiles/supervisors/supervisord.conf
[supervisord]
nodaemon = true
logfile = /var/log/supervisord.log
logfile_maxbytes = 10MB
pidfile = /var/run/supervisord.pid

[program:apache2]
command = apachectl -D "FOREGROUND"
autostart = true
autorestart = true

[program:cron]
command = cron -f
autostart = true
autorestart = true

[program:workers]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/rk-web/artisan queue:work redis --sleep=3 --tries=1 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=8
redirect_stderr=true
stdout_logfile=/var/log/supervisor/worker.log
stopwaitsecs=3600

#[program:websockets]
#command=php /var/www/rk-web/artisan websockets:serve
#numprocs=1
#autostart=true
#autorestart=true
